# telemetry/includes/cpu-telemetry.service
[Unit]
Description=Simple bash script

[Service]
User=robot
ExecStart=/usr/bin/python3 path_to_cpu_telemetry

[Install]
WantedBy=multi-user.target